# deploy/gigagochi-backup.timer
[Unit]
Description=Run Gigagochi off-site backup every night

[Timer]
OnCalendar=*-*-* 03:30:00 Europe/Moscow
Persistent=true
RandomizedDelaySec=15m
Unit=gigagochi-backup.service

[Install]
WantedBy=timers.target
